# php74/up-worker.conf
[program:up-worker]
command=php /opt/up-worker.php
autostart=true
autorestart=true
redirect_stderr=true
stderr_logfile=/application/worker.err.log
# stdout_logfile=/application/worker.out.log