# php80/worker.conf
[program:up-worker]
command=php $WORKER_PATH/$WORKER_FILE
autostart=true
autorestart=true
redirect_stderr=true
stderr_logfile=$WORKER_PATH/worker.err.log
stdout_logfile=$WORKER_PATH/worker.out.log