# php81/worker.conf
[program:up-worker]
command=php /opt/worker.php
autostart=true
autorestart=true
redirect_stderr=true
stderr_logfile=/opt/worker.err.log
stdout_logfile=/opt/worker.out.log